# etc/init.d/bootmisc
#!/sbin/openrc-run

description="Miscellaneous tasks to be done during bootup."

# Define dependencies
depend() {
    need localmount
    after mountnfs-bootclean
}

# Environment variables with default values
: ${DELAYLOGIN:=yes}

start_pre() {
    # If DELAYLOGIN is enabled, create the /run/nologin file
    if echo "${DELAYLOGIN}" | grep -iq "^y"; then
        echo "System bootup in progress - please wait" > /run/nologin || eerror "Failed to create /run/nologin"
    fi
}

start() {
    # Remove bootclean's flag files silently
    rm -f /tmp/.clean /run/.clean /run/lock/.clean
    rm -f /tmp/.tmpfs /run/.tmpfs /run/lock/.tmpfs

    # Reset /var/run/utmp
    local utmp='/var/run/utmp'
    if : > "${utmp}"; then
        chgrp utmp "${utmp}" || ewarn "Failed to chgrp ${utmp}"
        chmod 664 "${utmp}" || ewarn "Failed to chmod ${utmp}"
        [ -x /sbin/restorecon ] && /sbin/restorecon "${utmp}" || ewarn "Failed to restorecon ${utmp}"
    else
        eerror "Failed to truncate ${utmp}"
        return 1
    fi
}

stop() {
    # If DELAYLOGIN is enabled, remove the /run/nologin file
    if echo "${DELAYLOGIN}" | grep -iq "^y"; then
        rm -f /run/nologin || ewarn "Failed to remove /run/nologin"
    fi
}

status() {
    # Return status based on the presence of /run/nologin
    if [ -f /run/nologin ]; then
        return 0
    else
        return 1
    fi
}
