# etc/init.d/brightness
#!/sbin/openrc-run

description="Save and restore brightness level between restarts"

depend() {
    after mountall
}

SAVEFILE_PREFIX=/var/lib/initscripts/brightness

do_action() {
    local action=$1 label=$2 knob=$3 max=$4 file=$SAVEFILE_PREFIX${label:+.$label}
    [ -e "$knob" ] || return 0

    case "$action" in
        start)
            [ -f "$file" ] && cat "$file" >"$knob" || cat "$max" >"$knob"
            ;;
        stop)
            [ "$(cat "$knob")" -eq 0 ] && return 0
            cat "$knob" >"$file"
            ;;
        status)
            local msg="Current${label:+ $label} brightness level is $(cat "$knob")"
            if [ -f "$file" ]; then
                einfo "${msg}, saved value is $(cat "$file")"
                return 0
            else
                ewarn "${msg}, there is no saved value"
                return 4
            fi
            ;;
    esac
    return $?
}

start() {
    ebegin "Restoring brightness levels"
    do_action start '' /sys/class/backlight/acpi_video0/brightness /sys/class/backlight/acpi_video0/max_brightness
    do_action start intel /sys/class/backlight/intel_backlight/brightness /sys/class/backlight/intel_backlight/max_brightness
    eend $?
}

stop() {
    ebegin "Saving brightness levels"
    do_action stop '' /sys/class/backlight/acpi_video0/brightness /sys/class/backlight/acpi_video0/max_brightness
    do_action stop intel /sys/class/backlight/intel_backlight/brightness /sys/class/backlight/intel_backlight/max_brightness
    eend $?
}

status() {
    do_action status '' /sys/class/backlight/acpi_video0/brightness /sys/class/backlight/acpi_video0/max_brightness
    do_action status intel /sys/class/backlight/intel_backlight/brightness /sys/class/backlight/intel_backlight/max_brightness
}
