# etc/init.d/checkroot
#!/sbin/openrc-run

description="Check the root file system and manage root filesystem state"

depend() {
    before mountall
    after mountdevsubfs hostname
}

start() {
    ebegin "Checking root file system"
    if [ -f /forcefsck ] || grep -q -s -w -i "forcefsck" /proc/cmdline; then
        force="-f"
    else
        force=""
    fi

    if [ "$(findmnt -n -o FSTYPE /)" = "btrfs" ]; then
        ewarn "btrfs root detected, skipping root file system check."
    else
        fsck -T -a $force / || fsck -T -p $force /
        if [ $? -gt 3 ]; then
            eerror "Root file system check failed. A manual fsck is required."
            sulogin --force || {
                eerror "Failed to start maintenance shell. Rebooting in 5 seconds."
                sleep 5
                reboot -f
            }
        fi
    fi

    mount -o remount,rw / || mount -o remount,rw / || {
        eerror "Failed to remount root file system as read-write."
        return 1
    }

    if [ ! -h /etc/mtab ]; then
        rm -f /etc/mtab && ln -s /proc/mounts /etc/mtab
    fi

    eend 0
}

status() {
    rootrw=false
    swapon=false
    if grep -q ' / ' /proc/mounts | grep -q rw; then
        rootrw=true
    fi
    if [ -f /proc/swaps ] && grep -q -v '^Filename' /proc/swaps; then
        swapon=true
    fi
    if $rootrw || $swapon; then
        return 0
    else
        return 4
    fi
}
