# etc/init.d/checkroot-bootclean
#!/sbin/openrc-run

description="Clean temporary filesystems after the root filesystem is mounted"

depend() {
    after checkroot
    before bootmisc
}

clean_all() {
    which find >/dev/null 2>&1 || return 0
    ebegin "Cleaning up temporary files"
    local ES=0
    clean_tmp || ES=1
    clean /run "! -xtype d ! -name utmp ! -name innd.pid" || ES=1
    clean /run/lock "! -type d" || ES=1
    clean /run/shm "! -type d" || ES=1
    eend $ES
    return $ES
}

start() {
    ebegin "Cleaning temporary filesystems"
    rm -f /tmp/.clean /run/.clean /run/lock/.clean
    clean_all
    eend $?
}
