# etc/init.d/checkroot-bootclean
#!/sbin/openrc-run

description="Clean temporary filesystems after the root filesystem is mounted"

depend() {
    after checkroot
    before bootmisc
}

start() {
    ebegin "Cleaning temporary filesystems"
    rm -f /tmp/.clean /run/.clean /run/lock/.clean
    /lib/init/bootclean.sh
    eend $?
}
