# etc/init.d/hwclock
#!/sbin/openrc-run

description="Save system clock to hardware clock on shutdown"

depend() {
    after mountdevsubfs
    before umountfs
}

HCTOSYS_DEVICE="rtc0"

start() {
    ebegin "Skipping hwclock start as it is handled by udev rules"
    eend 0
}

stop() {
    ebegin "Saving the system clock to the hardware clock"
    if [ "$(uname -s)" = "GNU" ]; then
        /sbin/hwclock --directisa --systohc
    else
        /sbin/hwclock --rtc="/dev/$HCTOSYS_DEVICE" --systohc
    fi
    eend $?
}

reload() {
    stop
}

show() {
    ebegin "Displaying hardware clock time"
    if [ "$(uname -s)" = "GNU" ]; then
        /sbin/hwclock --directisa --show
    else
        /sbin/hwclock --rtc="/dev/$HCTOSYS_DEVICE" --show
    fi
    eend $?
}
