# etc/init.d/mountdevsubfs
#!/sbin/openrc-run

description="Mount special file systems under /dev"

depend() {
    need mountkernfs
    after udev
    before keyboard-setup
}

start() {
    ebegin "Mounting special filesystems under /dev"
    mount_filesystems mount_noupdate
    eend $?
}

stop() {
    # Typically, these filesystems are not unmounted on stop
    return 0
}

mount_filesystems() {
    MNTMODE="$1"

    # Ensure standard I/O nodes are present
    trylink() {
        test -h "/dev/$1" || test -d "/dev/$1" || ln -s -- "$2" "/dev/$1" || {
            eerror "Unable to ensure link /dev/$1 -> $2"
            ls -lad "/dev/$1"
        } >&2
    }

    trylink fd /proc/self/fd
    trylink stdin fd/0
    trylink stdout fd/1
    trylink stderr fd/2

    # Mount /dev/pts
    if [ "$(uname -s)" = "Linux" ]; then
        [ ! -d /dev/pts ] && mkdir --mode=755 /dev/pts
        [ -x /sbin/restorecon ] && /sbin/restorecon /dev/pts
        domount "$MNTMODE" devpts "" /dev/pts devpts "-onoexec,nosuid,gid=$TTYGRP,mode=$TTYMODE"
    fi
}

# Read defaults for TTYGRP and TTYMODE
[ -f /etc/default/devpts ] && . /etc/default/devpts
