# etc/init.d/mountall-bootclean
#!/sbin/openrc-run

description="Clean temporary filesystems after all local filesystems have been mounted."

depend() {
    need localmount
    before bootmisc
}

start() {
    ebegin "Cleaning /tmp, /var/lock, and /var/run"
    clean_temp_dirs
    eend $?
}

stop() {
    # No operation needed on stop
    return 0
}

clean_temp_dirs() {
    # Clean /tmp, /var/lock, /var/run
    # This script assumes you have proper mechanisms to handle the cleaning process.
    # Replace the following with whatever specific cleaning actions are required.
    find /tmp /var/lock /var/run -mindepth 1 -delete
}

