# etc/init.d/mountnfs
#!/sbin/openrc-run

description="Check and wait for network file systems to be mounted"

depend() {
    need mountall networking
}

start() {
    ebegin "Checking for network file systems"

    local waitnfs=""
    local ASYNCMOUNTNFS="${ASYNCMOUNTNFS:-yes}"
    local TIMEOUT=900

    for file in /etc/fstab; do
        [ -f "$file" ] || continue
        while read -r DEV MTPT FSTYPE OPTS REST; do
            case "$DEV" in
                ""|\#*) continue ;;
            esac
            case "$OPTS" in
                noauto|*,noauto|noauto,*|*,noauto,*) continue ;;
            esac
            case "$FSTYPE" in
                nfs|nfs4|smbfs|cifs|coda|ncp|ncpfs|ceph) ;;
                *) continue ;;
            esac
            case "$MTPT" in
                /usr/local|/usr/local/*|/usr|/usr/*|/var|/var/*)
                    waitnfs="$waitnfs $MTPT"
                    ;;
            esac
        done < "$file"
    done

    if [ "$ASYNCMOUNTNFS" != "no" ]; then
        for mountpt in $waitnfs; do
            ebegin "Waiting for $mountpt"
            while ! mountpoint -q "$mountpt"; do
                sleep 0.1
                TIMEOUT=$((TIMEOUT - 1))
                if [ "$TIMEOUT" -le 0 ]; then
                    ewarn "Timeout reached while waiting for $mountpt to be mounted."
                    eend 1
                    break
                fi
            done
            eend 0
        done
    else
        ewarn "ASYNCMOUNTNFS is set to 'no'. Skipping wait for network file systems."
    fi

    eend 0
}
