# etc/init.d/mountnfs
#!/sbin/openrc-run

description="Check for network file systems and log warnings if not mounted"

depend() {
    need localmount networking
}

start() {
    ebegin "Checking for network file systems"

    local waitnfs=""
    for file in /etc/fstab; do
        [ -f "$file" ] || continue
        while read -r DEV MTPT FSTYPE OPTS REST; do
            case "$DEV" in
                ""|\#*) continue ;;
            esac
            case "$OPTS" in
                noauto|*,noauto|noauto,*|*,noauto,*) continue ;;
            esac
            case "$FSTYPE" in
                nfs|nfs4|smbfs|cifs|coda|ncp|ncpfs|ceph) ;;
                *) continue ;;
            esac
            case "$MTPT" in
                /usr/local|/usr/local/*|/usr|/usr/*|/var|/var/*)
                    waitnfs="$waitnfs $MTPT"
                    ;;
            esac
        done < "$file"
    done

    for mountpt in $waitnfs; do
        if ! mountpoint -q "$mountpt"; then
            ewarn "Network file system $mountpt is not mounted. Boot process continuing."
        fi
    done

    eend 0
}
