# etc/init.d/mountall
#!/sbin/openrc-run

description="Mount all filesystems as specified in /etc/fstab"

depend() {
    need checkfs checkroot-bootclean
    before bootmisc hostname
}

: ${VERBOSE:=no}
: ${NOSWAP:=no}

mount_all_local() {
    if mountpoint -q /usr; then
        mount -o remount /usr || ewarn "Failed to remount /usr"
    fi

    mount -a -t noafs,nofs,norpc_pipefs,nonfs,nfs4,nonfs4,smbfs,cifs,ncp,ncpfs,coda,ocfs2,gfs,gfs2,ceph -O no_netdev || ewarn "Failed to mount some file systems"
}

remount_tmpfs() {
    local mountpoint="$1"
    if mountpoint -q "${mountpoint}"; then
        mount -o remount "${mountpoint}" || ewarn "Failed to remount ${mountpoint}"
    fi
}

start() {
    mount_all_local

    if [ ! -d /var/log/fsck ]; then
        mkdir -p /var/log/fsck || ewarn "Failed to create /var/log/fsck"
        chmod 755 /var/log/fsck || ewarn "Failed to set permissions on /var/log/fsck"
    fi

    if [ "$NOSWAP" != "yes" ]; then
        if [ "$VERBOSE" = "yes" ]; then
            swapon -a -e -v || ewarn "Failed to activate swap devices"
        else
            swapon -a -e 2>/dev/null || :
        fi
    else
        [ "$VERBOSE" = "no" ] || ewarn "Swap activation skipped due to 'noswap' boot parameter"
    fi

    remount_tmpfs "/run"
    remount_tmpfs "/run/lock"
    remount_tmpfs "/dev/shm"
    remount_tmpfs "/tmp"

    return 0
}

stop() {
    return 0
}

status() {
    if [ -d /var/log/fsck ]; then
        return 0
    else
        return 1
    fi
}
