# etc/init.d/mount-configfs
#!/sbin/openrc-run

description="Mount configfs kernel virtual file system"

depend() {
    after mountkernfs
    before kmod
}

start() {
    ebegin "Mounting configfs kernel virtual file system"
    local mountpoint="/sys/kernel/config"
    if [ -d "${mountpoint}" ] && ! findmnt "${mountpoint}" >/dev/null; then
        mount none -t configfs "${mountpoint}" || eend 1
    fi
    eend 0
}
