# etc/init.d/mount-configfs
#!/sbin/openrc-run

description="Mount configfs kernel virtual file system if the 'configfs' module is loaded."

depend() {
    after localmount
}

start() {
    ebegin "Mounting configfs kernel virtual file system"
    local mountpoint="/sys/kernel/config"

    # Ensure /sys/kernel/config exists and configfs is not already mounted
    if [ -d "${mountpoint}" ] && ! mountpoint -q "${mountpoint}"; then
        mount -t configfs configfs "${mountpoint}"
        eend $?
    else
        eend 0 "configfs is already mounted or not available."
    fi
}

stop() {
    # No specific stop operations needed; unmounting configfs is generally unnecessary
    return 0
}
