# etc/init.d/mountdevsubfs
#!/sbin/openrc-run

description="Mount special file systems under /dev"

depend() {
    after mountkernfs
    before eudev
}

start() {
    ebegin "Mounting special file systems under /dev"

    [ ! -h /dev/fd ] && [ ! -d /dev/fd ] && ln -s /proc/self/fd /dev/fd
    [ ! -h /dev/stdin ] && ln -s fd/0 /dev/stdin
    [ ! -h /dev/stdout ] && ln -s fd/1 /dev/stdout
    [ ! -h /dev/stderr ] && ln -s fd/2 /dev/stderr

    mount -t tmpfs tmpfs /run/shm || return 1

    if [ "$(uname -s)" = "Linux" ]; then
        [ ! -d /dev/pts ] && mkdir -m 755 /dev/pts
        mount -t devpts -o noexec,nosuid,gid=5,mode=620 devpts /dev/pts || return 1
    fi

    eend 0
}
