# etc/init.d/mountkernfs
#!/sbin/openrc-run

description="Mount kernel virtual file systems"

depend() {
    before mountall
}

start() {
    ebegin "Mounting kernel virtual file systems"
    mount -t tmpfs tmpfs /run || return 1
    mount -t tmpfs tmpfs /run/lock || return 1
    mount -t proc -o nodev,noexec,nosuid proc /proc || return 1
    if grep -E -qs "securityfs\$" /proc/filesystems; then
        mount -t securityfs securityfs /sys/kernel/security || return 1
    fi
    if grep -E -qs "sysfs\$" /proc/filesystems; then
        mount -t sysfs -o nodev,noexec,nosuid sysfs /sys || return 1
    fi
    if [ -d /sys/fs/pstore ]; then
        mount -t pstore pstore /sys/fs/pstore || return 1
    fi
    efivarsmnt=/sys/firmware/efi/efivars
    if [ -d "$efivarsmnt" ]; then
        mount -t efivarfs none "$efivarsmnt" || return 1
    fi
    eend 0
}
