# etc/init.d/umountfs
#!/sbin/openrc-run

description="Turn off swap and unmount all local file systems"

depend() {
    after umountroot
}

stop() {
    ebegin "Unmounting file systems and deactivating swap"

    local PROTECTED_MOUNTS REG_MTPTS TMPFS_MTPTS WEAK_MTPTS
    PROTECTED_MOUNTS=$(sed -n ':a;/^[^ ]* \/ /!{H;n;ba};{H;s/.*//;x;s/\n//;p}' /proc/mounts)

    while read -r DEV MTPT FSTYPE REST; do
        echo "$PROTECTED_MOUNTS" | grep -qsF "$DEV $MTPT " && continue
        case "$MTPT" in
            /|/usr|/proc|/dev|/.dev|/dev/pts|/dev/shm|/dev/.static/dev|/proc/*|/sys|/sys/*|/run|/run/lock|/run/shm|/run/rpc_pipefs|/dev/vcs) continue ;;
        esac
        case "$FSTYPE" in
            proc|procfs|linprocfs|sysfs|securityfs|usbfs|usbdevfs|devpts) continue ;;
            tmpfs) TMPFS_MTPTS="$MTPT $TMPFS_MTPTS" ;;
            *) 
                if echo "$PROTECTED_MOUNTS" | grep -qs "^$DEV "; then
                    WEAK_MTPTS="$MTPT $WEAK_MTPTS"
                else
                    REG_MTPTS="$MTPT $REG_MTPTS"
                fi
                ;;
        esac
    done < /proc/mounts

    if [ "$TMPFS_MTPTS" ]; then
        einfo "Unmounting temporary filesystems"
        umount $TMPFS_MTPTS
    fi

    einfo "Deactivating swap"
    swapoff -a

    if [ "$WEAK_MTPTS" ]; then
        einfo "Unmounting weak filesystems"
        umount -r -d $WEAK_MTPTS
    fi

    if [ "$REG_MTPTS" ]; then
        einfo "Unmounting local filesystems"
        umount -f -r -d $REG_MTPTS
    fi

    eend 0
}
