# etc/init.d/urandom
#!/sbin/openrc-run

description="Save and restore random seed between restarts"

depend() {
    after mountall
}

SAVEDFILE=/var/lib/urandom/random-seed
POOLBYTES=$(($(cat /proc/sys/kernel/random/poolsize 2>/dev/null || echo 4096) / 8))

start() {
    ebegin "Initializing random number generator"
    (
        date +%s.%N
        [ -f "$SAVEDFILE" ] && cat "$SAVEDFILE"
    ) >/dev/urandom

    umask 077
    dd if=/dev/urandom of="$SAVEDFILE" bs="$POOLBYTES" count=1 >/dev/null 2>&1
    eend $?
}

stop() {
    ebegin "Saving random seed"
    umask 077
    dd if=/dev/urandom of="$SAVEDFILE" bs="$POOLBYTES" count=1 >/dev/null 2>&1
    eend $?
}

status() {
    [ -f "$SAVEDFILE" ]
}
